# src/supervisord.conf
[supervisord]
nodaemon = true

[program:mysqld]
command = /usr/bin/mysql_start.sh

[program:nginx]
command = /usr/sbin/nginx

[program:php-fpm7.0]
command = /usr/sbin/php-fpm7.0 -F
